# Config/vsido.streamingface.service
[Unit]
Description=VSido Stream Face Detect Video Service
After=syslog.target network.target home.mount
	
[Service]
Type=simple
ExecStart= /home/sysroot/usr/bin/MJPEStream.face.sh >/dev/null
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target

